# service/runn.service
[Unit]
Description=BadVPN UDPGW (UDP forwarder for SSH/SSL clients)
After=network.target

[Service]
Type=simple
ExecStart=/usr/local/bin/run-udpgw
Restart=on-failure
RestartSec=3
KillMode=control-group

[Install]
WantedBy=multi-user.target
